# chaosmap.service
[Unit]
Description=Chaos Zone Map Flask Application
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/path/to/your/chaos-zone-map
Environment=PATH=/path/to/your/chaos-zone-map/bin
ExecStart=/path/to/your/chaos-zone-map/bin/python -m flask run --host=127.0.0.1 --port=5000
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
